# backup/crontab/daily-backup
# Daily backup at midnight run backup
0 0 * * * /bin/backup >> /dev/null 2>&1

# Daily at 00:30 prune old backups
30 0 * * * /bin/prune >> /dev/null 2>&1
